# config/systemd/mywebapp.socket
[Unit]
Description=Socket for MyWebApp

[Socket]
ListenStream=5200

[Install]
WantedBy=sockets.target