# config/systemd/mywebapp.service
[Unit]
Description=MyWebApp Service
After=network.target postgresql.service

[Service]
Type=simple
User=app
Group=app
WorkingDirectory=/opt/mywebapp
Environment=CONFIG_PATH=/etc/mywebapp/config.json
ExecStartPre=/usr/bin/node /opt/mywebapp/src/db/migrate.js
ExecStart=/usr/bin/node /opt/mywebapp/src/app.js
Restart=always

[Install]
WantedBy=multi-user.target